# debian/pomodoro-lock/usr/share/pomodoro-lock/config/pomodoro-lock.service
[Unit]
Description=Pomodoro Lock UI - Standalone Timer Application
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
WorkingDirectory=/usr/share/pomodoro-lock
ExecStart=/usr/bin/pomodoro-lock ui
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical-session.target 