# pomodoro-lock.service
[Unit]
Description=Pomodoro Lock Service
After=graphical.target
Wants=graphical.target

[Service]
Type=simple
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/vinay/.Xauthority
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus
ExecStartPre=/bin/sleep 10
ExecStart=/usr/bin/python3 /home/vinay/.local/share/pomodoro-lock/bin/pomodoro-lock.py
Restart=always
RestartSec=5

[Install]
WantedBy=graphical.target 