# debian/postinst
#!/bin/bash
# Post-installation script for pomodoro-lock

set -e

case "$1" in
    configure)
        echo "Pomodoro Lock - Standalone UI Installation"
        echo "=========================================="
        echo ""
        echo "The application has been installed to /usr/share/pomodoro-lock"
        echo "User data will be stored in ~/.local/share/pomodoro-lock"
        echo ""
        echo "New Architecture:"
        echo "  - Standalone UI: Complete timer application with integrated overlays"
        echo "  - Systemd Service: Manages autostart and process lifecycle"
        echo ""
        echo "Available commands:"
        echo "  pomodoro-lock      - Start UI (default)"
        echo "  pomodoro-lock ui   - Start UI"
        echo "  pomodoro-lock service - Start the systemd service"
        echo "  pomodoro-lock stop - Stop the systemd service"
        echo "  pomodoro-lock status - Show service status"
        echo "  pomodoro-configure - Configure the application"
        echo ""
        echo "Recommended usage:"
        echo "  1. Start the UI: pomodoro-lock ui"
        echo "  2. Or enable autostart: systemctl --user enable pomodoro-lock.service"
        echo ""
        echo "The first time you run 'pomodoro-lock', it will:"
        echo "  1. Create a virtual environment in ~/.local/share/pomodoro-lock/venv"
        echo "  2. Install Python dependencies (including notification support)"
        echo "  3. Set up user configuration"
        echo "  4. Configure systemd user service"
        echo ""
        echo "Features:"
        echo "  - Desktop notifications before break periods"
        echo "  - Multi-display screen overlays during breaks"
        echo "  - System tray integration"
        echo "  - Configurable work/break cycles"
        echo ""
        echo "To get started, run: pomodoro-lock"
        echo ""
        ;;
    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

exit 0 