# config/pomodoro-lock.service
[Unit]
Description=Pomodoro Lock Service with Multi-Display Overlay
After=graphical-session.target

[Service]
Type=simple
Restart=on-failure
RestartSec=2
Environment=PYTHONUNBUFFERED=1
WorkingDirectory=%h/.local/share/pomodoro-lock
ExecStart=%h/.local/share/pomodoro-lock/venv/bin/python3 %h/.local/share/pomodoro-lock/pomodoro-service.py

[Install]
WantedBy=graphical-session.target 