# config/pomodoro-lock-simple.service
[Unit]
Description=Pomodoro Lock Service
After=display-manager.service
Wants=display-manager.service

[Service]
Type=simple
User=vinay
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/vinay/.Xauthority
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus
Environment=XDG_RUNTIME_DIR=/run/user/1000
WorkingDirectory=/home/vinay/.local/share/pomodoro-lock
ExecStartPre=/bin/bash -c 'until xset q >/dev/null 2>&1; do sleep 1; done'
ExecStartPre=/bin/sleep 10
ExecStart=/usr/bin/python3 /home/vinay/.local/share/pomodoro-lock/bin/pomodoro-lock.py
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
TimeoutStartSec=60

[Install]
WantedBy=default.target 